# rtorrent-flood/rootfs/etc/services.d/rtorrent/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# rTorrent
# ==============================================================================
bashio::log.info "Starting rTorrent..."

exec /usr/bin/rtorrent.sh
