# rtorrent-flood/rootfs/etc/services.d/flood/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Flood
# ==============================================================================
bashio::log.info "Starting Flood UI..."
exec /usr/bin/flood.sh
